# systemd/pi-deployer.service
[Unit]
Description=Pi-Deployer Webhook Service
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=pie
Group=pie
WorkingDirectory=/home/pie/pi-deployer
EnvironmentFile=/home/pie/pi-deployer/.env
ExecStart=/home/pie/pi-deployer/venv/bin/python /home/pie/pi-deployer/deployer.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-deployer

[Install]
WantedBy=multi-user.target
