# systemd/pi-deployer.service
[Unit]
Description=Pi-Deployer Webhook Service
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pi-deployer
EnvironmentFile=/home/pi/pi-deployer/.env
ExecStart=/usr/bin/python3 /home/pi/pi-deployer/deployer.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pi-deployer

[Install]
WantedBy=multi-user.target
